# all-in-one/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:klocwork]
command=/home/klocwork/server/bin/kwservice -r /home/klocwork/server/projects_root start